# deploy/usb-cam-host/com.farmguardian.usb-cam-host.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  usb-cam-host LaunchAgent. Runs as the logged-in user so AVFoundation can
  surface a Camera TCC prompt on first launch (LaunchDaemons run as root,
  cannot surface a user-visible TCC prompt, and cannot be granted user-scope
  Camera access on modern macOS).

  Install:
    cp com.farmguardian.usb-cam-host.plist ~/Library/LaunchAgents/
    launchctl unload ~/Library/LaunchAgents/com.farmguardian.usb-cam-host.plist 2>/dev/null
    launchctl load  ~/Library/LaunchAgents/com.farmguardian.usb-cam-host.plist

  Tailor per-host:
    - ProgramArguments: point at the Python venv and the repo checkout on that box
    - USB_CAM_PORT: keep 8089 unless it collides
    - USB_CAM_DEVICE_INDEX: verify per-host — do NOT assume 0=built-in, 1=USB.
        On MacBook Air (192.168.0.50): index 0 = external USB webcam (enumerates
        first), index 1 = built-in FaceTime HD. To verify on any host:
        python3 -c "import cv2; [print(i, cv2.VideoCapture(i).isOpened()) for i in range(4)]"
        then cross-check with: system_profiler SPCameraDataType
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.farmguardian.usb-cam-host</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/macmini/Documents/GitHub/farm-guardian/venv/bin/python</string>
        <string>/Users/macmini/Documents/GitHub/farm-guardian/tools/usb-cam-host/usb_cam_host.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/macmini/Documents/GitHub/farm-guardian</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>USB_CAM_PORT</key>
        <string>8089</string>
        <key>USB_CAM_DEVICE_INDEX</key>
        <string>0</string>
        <key>USB_CAM_WIDTH</key>
        <string>1920</string>
        <key>USB_CAM_HEIGHT</key>
        <string>1080</string>
        <key>USB_CAM_WARMUP</key>
        <string>15</string>
        <key>USB_CAM_JPEG_QUALITY</key>
        <string>95</string>
        <key>USB_CAM_AUTO_WB</key>
        <string>true</string>
        <key>USB_CAM_WB_STRENGTH</key>
        <string>0.5</string>

        <!--
          Manual exposure — starts at -6 (≈1/64 sec on AVFoundation). The
          generic UVC webcam's auto-exposure clips the red channel to 255
          under the brooder heat lamp, which is what causes the "nuclear
          pink" look after gray-world WB runs. Clamp the shutter below the
          clipping point so the WB pipeline has real color data to work
          with. Tune by watching a pure-red pixel crop: if any channel
          still reads 255, drop one more stop (-7). If the frame is too
          dark and the chicks become silhouettes, step back up (-5).
          See docs/16-Apr-2026-heat-lamp-orange-cast-investigation.md.
        -->
        <key>USB_CAM_AUTO_EXPOSURE</key>
        <string>manual</string>
        <key>USB_CAM_EXPOSURE</key>
        <string>-6</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>OPENCV_AVFOUNDATION_SKIP_AUTH</key>
        <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/usb-cam-host.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/usb-cam-host.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
