# deploy/s7-settings-watchdog/com.farmguardian.s7-settings-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.farmguardian.s7-settings-watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>ts=$(date -u +%Y-%m-%dT%H:%M:%SZ); wb_code=$(/usr/bin/curl -sS -m 5 -o /dev/null -w '%{http_code}' 'http://192.168.0.249:8080/settings/whitebalance?set=incandescent'); fm_code=$(/usr/bin/curl -sS -m 5 -o /dev/null -w '%{http_code}' 'http://192.168.0.249:8080/settings/focusmode?set=continuous-picture'); echo "$ts wb=$wb_code fm=$fm_code" &gt;&gt; /tmp/s7-settings-watchdog.log</string>
    </array>

    <key>StartInterval</key>
    <integer>600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/s7-settings-watchdog.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/s7-settings-watchdog.err.log</string>
</dict>
</plist>
